# etc/supervisord-development.conf
[inet_http_server]
port=127.0.0.1:8988

[supervisord]
logfile=%(here)s/../var/supervisord.log
logfile_maxbytes=100MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/../var/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:8988

[program:zeo]
command = %(here)s/../bin/runzeo -C %(here)s/zeo.conf
redirect_stderr = true
stdout_logfile = %(here)s/../var/zeo.log
autostart = true
stopwaitsecs = 30
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:yss]
command = %(here)s/../bin/pserve %(here)s/development.ini
redirect_stderr = true
stdout_logfile = %(here)s/../var/yss.log
autostart = true
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:indexer]
command = %(here)s/../bin/sd_drain_indexing %(here)s/development.ini
redirect_stderr = true
stdout_logfile = %(here)s/../var/indexing.log
autostart = true
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:postproc]
command = %(here)s/../bin/postproc %(here)s/development.ini
redirect_stderr = true
stdout_logfile = %(here)s/../var/postproc.log
autostart = true
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:retimer]
command = %(here)s/../bin/retimer %(here)s/development.ini
redirect_stderr = true
stdout_logfile = %(here)s/../var/retime.log
autostart = true
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[program:sd_drain_indexing]
command = %(here)s/../bin/sd_drain_indexing %(here)s/development.ini
redirect_stderr = true
stdout_logfile = %(here)s/../var/drain_indexing.log
autostart = true
startsecs = 5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
